# update-motd.d/00-0logo
#!/bin/sh
# To add the GSL logo to the login message
#    1. Install the update-motd package
#    2. Install this file as /etc/update-motd.d/00-0logo
#    3. Copy the file ../gsl-motd-logo.txt to /etc
[ -f /etc/gsl-motd-logo.txt ] && cat /etc/gsl-motd-logo.txt
