# profile.d/drush.sh
# Set the default PHP version under which drush runs. Before you change this
# globally, consider whether setting an environment variable for yourself is
# a more appropriate solution.
export DRUSH_PHP='/usr/bin/php5.6'
